# conf/systemd.service
[Unit]
Description=TheLounge: IRC client
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__INSTALL_DIR__/
Environment="PATH=__ENV_PATH__"
Environment="THELOUNGE_HOME=__DATA_DIR__/"
Environment="NODE_ENV=production"
ExecStart=/usr/bin/yarn start
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=multi-user.target
