# conf/systemd.service
[Unit]
Description=The Lounge IRC client
After=__APP__.service

[Service]
Type=simple
ExecStart=/usr/bin/thelounge start --home /home/yunohost.app/__APP__/
User=__APP__
Group=__APP__
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=default.target

