# conf/systemd.service
[Unit]
Description=The Lounge IRC client
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__FINALPATH__
Environment="PATH=__ENV_PATH__"
Environment="THELOUNGE_HOME=/home/yunohost.app/__APP__"
Environment="NODE_ENV=production"
ExecStart=__YNH_NPM__ start -c port=__PORT__
Restart=on-failure
RestartSec=5
StartLimitInterval=60s
StartLimitBurst=3

[Install]
WantedBy=default.target
