# conf/systemd.service
[Unit]
Description=The Lounge IRC client
After=network.target

[Service]
Type=simple
User=__APP__
Group=__APP__
WorkingDirectory=__FINALPATH__
ExecStartPre=__NODEJS__
Environment="PATH=__ENV_PATH__"
ExecStart=__NODE__ index.js start --home /home/yunohost.app/__APP__/
Restart=always

[Install]
WantedBy=default.target

